# deployment/application/celery/celeryd.conf
[program:celeryd]
directory = /opt/Assassins/Assassins
command = /opt/.env/bin/celery worker -A Assassins -B --detach