# deployment/application/celery/flowerd.conf
[program:flowerd]
directory = /opt/Assassins/Assassins
command = /opt/.env/bin/celery -A Assassins flower